# docker/dockerfile/php/config/supervisord.conf
[supervisord]
nodaemon=true
loglevel = info
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ethos-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan queue:work --sleep=6 --tries=3
autostart=true
autorestart=true
numprocs=100
redirect_stderr=true
stdout_logfile=/var/log/schedule.log
stopwaitsecs=3600